# system/autostart/gym.desktop
[Desktop Entry]
Type=Application
Name=Gym
Exec=lxterminal -e "cd ~/Gym_Cards/;python3 gym.py;exec bash"
