# system/crontab
@reboot bash ~/Gym_Cards/system/startup.sh >>~/Gym_Cards/logs/system.log 3>&1
* * * * * bash ~/Gym_Cards/system/backuo.sh $(date +"%Y_%m_%d_%H") >>~/Gym_Cards/logs/system.log 2>&1